# extra/debian/hiawatha.postinst
#!/bin/sh

set -e
#DEBHELPER#
update-rc.d php-fcgi defaults > /dev/null
