# extra/debian/hiawatha.postrm
#!/bin/sh

set -e
#DEBHELPER#
update-rc.d -f php-fcgi remove >/dev/null || exit $?
